# dev/development/files/etc/cron.d/vagrant-ntpdate
PATH=/usr/bin:/usr/sbin:/bin:/sbin

1 * * * * root service ntp stop; ntpdate 10.0.2.2; service ntp start
